# proxmox-ups.service
[Unit]
Description=Foo

[Service]
User=root
ExecStart=/bin/bash -c "proxmox-ups"

[Install]
WantedBy=multi-user.target